# ops/systemd/openclaw-gateway@.service
# Per-user OpenClaw gateway systemd template unit.
# Enable with: systemctl enable openclaw-gateway@<linux-user>.service

[Unit]
Description=ClawaaS OpenClaw gateway for %i
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i
Environment=HOME=/home/%i
Environment=OPENCLAW_HOME=/home/%i/.openclaw

# cmdok reads config from $OPENCLAW_HOME/openclaw.json
# Port is baked into the per-user config by render_openclaw_config.py
ExecStart=/usr/local/bin/cmdok gateway start --config /home/%i/.openclaw/openclaw.json
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=15

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/%i

[Install]
WantedBy=multi-user.target
