# ops/systemd/openclaw-gateway@.service
# Per-user OpenClaw gateway systemd template unit.
# Enable with: systemctl enable openclaw-gateway@<linux-user>.service

[Unit]
Description=ClawaaS OpenClaw gateway for %i
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i
Environment=HOME=/home/%i
Environment=OPENCLAW_HOME=/home/%i
Environment=TMPDIR=/home/%i/.openclaw/tmp
Environment=NODE_OPTIONS=--use-system-ca
EnvironmentFile=/home/%i/.openclaw/gateway.env

# openclaw reads config from $OPENCLAW_HOME/openclaw.json automatically
# Port from gateway.env; --allow-unconfigured skips gateway.mode=local requirement
ExecStartPre=/bin/mkdir -p /home/%i/.openclaw/tmp
ExecStart=/usr/local/bin/openclaw gateway run --allow-unconfigured --port ${CLAWAAS_GATEWAY_PORT}
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=15

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/%i /tmp
ReadOnlyPaths=/etc/ssl /etc/resolv.conf /etc/hosts /usr/share/ca-certificates /opt/openclaw

[Install]
WantedBy=multi-user.target
